# files/etc/profile.d/openxr.sh
export XR_RUNTIME_JSON=/usr/share/openxr/1/openxr_monado.json
