# files/etc/profile.d/rclone.sh
export RCLONE_CONFIG=/etc/rclone.conf
export RCLONE_CONFIG_H_HOST=$s
